# templates/disable-thp.service.j2
# =============================================================================
# Disable Transparent Huge Pages Service
# =============================================================================
# Generated by Ansible - Do not edit manually
# Required for MongoDB performance
# =============================================================================

[Unit]
Description=Disable Transparent Huge Pages (THP)
DefaultDependencies=no
After=sysinit.target local-fs.target
Before=mongod.service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo never > /sys/kernel/mm/transparent_hugepage/enabled'
ExecStart=/bin/sh -c 'echo never > /sys/kernel/mm/transparent_hugepage/defrag'
RemainAfterExit=yes

[Install]
WantedBy=basic.target
