# templates/mongos.service.j2
# =============================================================================
# MongoDB Mongos Router Service
# =============================================================================
# Generated by Ansible - Do not edit manually
# =============================================================================

[Unit]
Description=MongoDB Sharding Router (mongos)
Documentation=https://docs.mongodb.com/manual/reference/program/mongos/
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User={{ mongodb_user }}
Group={{ mongodb_group }}
Environment="OPTIONS=-f /etc/mongos.conf"
ExecStart=/usr/bin/mongos $OPTIONS
PIDFile={{ mongodb_pid_directory }}/mongos.pid

# Resource limits
LimitFSIZE=infinity
LimitCPU=infinity
LimitAS=infinity
LimitNOFILE={{ mongodb_limit_nofile }}
LimitNPROC={{ mongodb_limit_nproc }}
LimitMEMLOCK={{ mongodb_limit_memlock }}

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
