# templates/mongodb-monitor.service.j2
[Unit]
Description=MongoDB Health Monitor
Documentation=https://docs.mongodb.com/

[Service]
Type=oneshot
ExecStart=/opt/mongodb-alerting/scripts/mongodb-health-monitor.sh
User=root
Group=root
