# templates/mongodb_exporter.service.j2
# =============================================================================
# MongoDB Prometheus Exporter Service
# =============================================================================
# Generated by Ansible - Do not edit manually
# =============================================================================

[Unit]
Description=MongoDB Prometheus Exporter
Documentation=https://github.com/percona/mongodb_exporter
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=simple
User=mongodb_exporter
Group=mongodb_exporter
EnvironmentFile=/etc/mongodb_exporter/mongodb_exporter.env
ExecStart=/usr/local/bin/mongodb_exporter \
    --web.listen-address=:{{ mongodb_exporter_port }} \
    --web.telemetry-path={{ mongodb_exporter_path }} \
    --collect-all \
    --compatible-mode

Restart=on-failure
RestartSec=5

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full

[Install]
WantedBy=multi-user.target
